# sgbus-bot.service
[Unit]
Description=SG Bus Arrival Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/sgbus-bot
ExecStart=/sgbus-bot/.venv/bin/python bot.py
Restart=always
RestartSec=10
KillMode=control-group
TimeoutStopSec=15
StandardOutput=append:/sgbus-bot/data/bot.log
StandardError=append:/sgbus-bot/data/bot.log

[Install]
WantedBy=multi-user.target
